# .config/systemd/user/darkman@.timer
[Unit]
Description=Switch to %i mode on schedule

[Timer]
Persistent=true

[Install]
WantedBy=timers.target
