# .config/systemd/user/darkman@.service
[Unit]
Description=Switch to %i mode

[Service]
Type=oneshot
ExecStart=/usr/bin/darkman set %i
