# .config/systemd/user/darkman-set.service
[Unit]
Description=Set darkman mode based on time

[Service]
Type=oneshot
ExecStart=%h/.local/bin/darkman-set
