# .config/systemd/user/obsidian-sync.timer
[Unit]
Description=Timer for obsidian

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m
