# .config/systemd/user/obsidian-sync.service
[Unit]
Description=Sync obsidian vault

[Service]
Type=oneshot
WorkingDirectory=%h/.obsidian-vault
ExecStart=git-sync
