# .config/systemd/user/redshift-hourly.service
[Unit]
Description=redshift
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pm redshift
