# .config/systemd/user/socks-proxy.socket
[Unit]
Description=SOCKS proxy

[Socket]
ListenStream=1080

[Install]
WantedBy=sockets.target
